# .devcontainer/config/mysql.conf
[program:mysql]
command=/usr/sbin/mysqld --user=mysql --bind-address=0.0.0.0
priority=1000
autostart=true
autorestart=true
user=root